# net-print/cloudprint/files/cloudprint
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-emulation/virtualbox/files/vboxwebsrv-initd,v 1.1 2011/01/06 22:47:34 polynomial-c Exp $

opts="setup logout" 

depend() {
	need net
}

start() {
	ebegin "Starting Cloudprint service"
	start-stop-daemon --start --quiet -m --pidfile /var/run/cloudprint.pid \
		--background -i --exec /usr/lib/cloudprint/cloudprint.py
	eend $?
}

stop() {
        ebegin "Stopping Cloudprint service"
        start-stop-daemon --stop --quiet --pidfile /var/run/cloudprint.pid
        eend $?
}

setup() {
	ebegin "Setting up Cloudprint service"
	/usr/lib/cloudprint/cloudprint.py -s
	eend $?
}

logout() {
	ebegin "Removing Cloudprint authentication details"
	/usr/lib/cloudprint/cloudprint.py -l
	eend $?
}
